# docker/supervisord.conf
[supervisord]
nodaemon=true
logfile=/dev/stdout
logfile_maxbytes=0
pidfile=/tmp/supervisord.pid

[program:portal-api]
command=node services/portal-api/dist/index.js
directory=/app
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=10

[program:indexer-supplier]
command=node services/indexer/dist/cli.js infra/configs/indexer-supplier.yaml --serve
directory=/app
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=20

[program:indexer-buyer]
command=node services/indexer/dist/cli.js infra/configs/indexer-buyer.yaml --serve
directory=/app
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=20

[program:indexer-financier-a]
command=node services/indexer/dist/cli.js infra/configs/indexer-financier-a.yaml --serve
directory=/app
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=20

[program:indexer-financier-b]
command=node services/indexer/dist/cli.js infra/configs/indexer-financier-b.yaml --serve
directory=/app
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=20

[program:indexer-regulator]
command=node services/indexer/dist/cli.js infra/configs/indexer-regulator.yaml --serve
directory=/app
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=20

[program:indexer-platform]
command=node services/indexer/dist/cli.js infra/configs/indexer-platform.yaml --serve
directory=/app
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=20

[program:notifications]
command=node services/notifications/dist/cli.js
directory=/app
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=30

[program:oracle-relay]
command=node services/oracle-relay/dist/cli.js
directory=/app
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=30

[program:registry-api]
command=node services/registry-api/dist/index.js
directory=/app
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=30

[program:kyb-gateway]
command=node services/kyb-gateway/dist/index.js
directory=/app
environment=KYB_DATA_DIR="/data/kyb",KYB_GATEWAY_PORT="8090",KYB_COMPLETE_SECRET="%(ENV_KYB_COMPLETE_SECRET)s"
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=40

[program:party-provisioner]
command=node services/party-provisioner/dist/index.js
directory=/app
environment=PROVISIONER_DATA_DIR="/data/provisioner",PROVISIONER_PORT="8091",KYB_GATEWAY_URL="http://127.0.0.1:8090"
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=40

[program:agent-runtime]
command=/app/docker/run-agent-runtime.sh
directory=/app
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=50

[program:nginx]
command=nginx -g "daemon off;"
directory=/app
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=100
